# nrf-udev_1.0.1-all/DEBIAN/postinst
#!/bin/sh -e

echo "Reloading udev rules..."
invoke-rc.d udev reload

